# blacklists/blacklists.service
[Unit]
Description=updating ip address blacklists

[Service]
Type=oneshot
ExecStart=/usr/bin/blacklists

