# blacklists/blacklists.service
[Unit]
Description=updating ip address blacklist

[Service]
Type=oneshot
ExecStart=/usr/bin/blacklists

