# blacklists/blacklists.service
[Unit]
Description=updating ntopng ipset from ntopng blacklists

[Service]
Type=oneshot
ExecStart=/usr/bin/blacklists

